# packages/karma-0.68/services/karma.service
[Unit]
Description=karma
Wants=network-online.target
After=network-online.target

[Service]
User=karma
Group=karma
Type=simple
WorkingDirectory=/etc/karma/
ExecStart=/usr/local/bin/karma --config.file=/etc/karma/karma.yml

[Install]
WantedBy=multi-user.target
