# packages/ipmi_exporter/services/ipmi_exporter.service
[Unit]
Description=ipmi_exporter
Wants=network-online.target
After=network-online.target

[Service]
User=ipmi_exporter
Group=ipmi_exporter
Type=simple
WorkingDirectory=/etc/ipmi_exporter/
ExecStart=/usr/local/bin/ipmi_exporter --config.file=/etc/ipmi_exporter/ipmi_config.yml

[Install]
WantedBy=multi-user.target
